# dot_config/autostart/polybar.desktop
[Desktop Entry]
Name=Polybar
Exec=/home/sam/.config/polybar/launch.sh
Terminal=false
Type=Application